# etc/development/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/
user=root
[program:php-fpm]
process_name=%(program_name)s
command=php-fpm --nodaemonize
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/php-fpm.log
priority=5
[program:nginx]
process_name=%(program_name)s
command= /bin/sh -c "sleep 3 && nginx -g 'daemon off;'"
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/nginx.log
priority=10